# powersmtp.init.d.in
#!/bin/sh

# chkconfig: - 80 75
# description: PowerSMTP receives incoming email

prefix=@prefix@
exec_prefix=@exec_prefix@
bindir=@bindir@
sbindir=@sbindir@

case "$1" in
	start | stop | status)
	$sbindir/powersmtp $1
	;;

	force-reload | restart)
	$sbindir/powersmtp stop
	$sbindir/powersmtp start
	;;

	monitor)
		$sbindir/powersmtp --daemon=no --loglevel=9
	;;
	
	*)
	echo powersmtp [start\|stop\|restart\|force-reload\|status]
	;;
esac

