# pplistener.init.d.in
#!/bin/sh

# chkconfig: - 70 75
# description: PowerSMTP pplistener manages access to disk storage,
#              making it available to the PowerSMTP and PowerPOP services.

prefix=@prefix@
exec_prefix=@exec_prefix@
bindir=@bindir@
sbindir=@sbindir@

case "$1" in
	start | stop | status)
	$sbindir/pplistener $1
	;;

	force-reload | restart)
	$sbindir/pplistener stop
	$sbindir/pplistener start
	;;
	
	*)
	echo pplistener [start\|stop\|restart\|force-reload\|status]
	;;
esac

